# roombyroom/Controller/pi.cron
* * * * * /home/pi/run
@reboot sleep 10 && /home/pi/init
